# packaging/macos/co.securityops.torando-gui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- SPDX-License-Identifier: AGPL-3.0-only -->
<!-- Copyright (c) 2026 Cristian Cezar Moisés — AGPL-3.0-only -->
<!-- LaunchDaemon for the Torando Control root daemon. Runs as root so it can
     drive pf, the system SOCKS proxy (networksetup) and interface DNS. Install
     to /Library/LaunchDaemons and load with launchctl. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.securityops.torando-gui</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/torando-guid</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/torando-gui.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/torando-gui.log</string>
</dict>
</plist>
